# utils/s6-overlay/services/postgres/run
#!/command/with-contenv bash
# ^ allows access to the Docker container's environment variables

# -a  Mark variables which are modified or created for export
# -e  Exit immediately if a command exits with a non-zero status
set -ea
   
# start postgres with specified data directory, copying stderr onto stdout
exec s6-setuidgid postgres /usr/lib/postgresql/13/bin/postgres -c random_page_cost=1.0 -c shared_preload_libraries=pg_stat_statements -D /datadir/postgresql/data 2>&1